# systemd/sybase-shutdown.target
[Unit]
Description=Shut down  all Sybase services

Before=shutdown.target

[Install]
WantedBy=shutdown.target

